# supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=python /app/app.py
autostart=true
autorestart=true
stdout_logfile=/var/log/app.log
stderr_logfile=/var/log/app_err.log
